# csf-backend.service
[Unit]
Description=CSF Backend FastAPI Application
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=notify
User=root
Group=root
WorkingDirectory=/root/csf_backend
Environment="PATH=/root/csf_backend/.venv/bin"
Environment="DATABASE_URL=postgresql+asyncpg://csf_user:your_password_here@localhost:5432/csf_db"
Environment="SECRET_KEY=your-secret-key-min-32-characters-long"
Environment="CORS_ALLOWED_ORIGINS=*"
Environment="APP_ENV=production"
Environment="STRIPE_SECRET_KEY=your_stripe_key"
Environment="STRIPE_PUBLISHABLE_KEY=your_stripe_pub_key"
Environment="STRIPE_WEBHOOK_SECRET=your_stripe_webhook_secret"
ExecStart=/root/csf_backend/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000 --workers 4
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
